# packages/3rdParty/service/aml-vnc/source/system.d/service.aml-vnc.service
[Unit]
Description=Amlogic VNC Server
After=network-online.target kodi.service
Requires=kodi.service

[Service]
ExecStartPre=/bin/sleep 5
ExecStart=/storage/.kodi/addons/service.aml-vnc/bin/aml-vnc.start
SyslogIdentifier=aml-vnc
TimeoutStopSec=1
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.service
Alias=aml-vnc.service
