# packages/3rdParty/service/makemkv/source/system.d/service.makemkv.service
[Unit]
Description=MakeMKV bootstrap service
After=network-online.target

[Service]
Type=oneshot
ExecStart=/storage/.kodi/addons/service.makemkv/bin/libmakemkv-bootstrap.start
SyslogIdentifier=makemkv
RemainAfterExit=yes

[Install]
WantedBy=kodi.target
Alias=makemkv.service
